# init_instavue.conf
[program:init]
directory=/app
command=/app/scripts/start.sh